# debian/postinst
#!/bin/sh

systemctl is-enabled named.service && systemctl reload-or-restart named.service
